# logstash-v2/2.3/source/my_init/22_logstash
#! /bin/sh

### BEGIN INIT INFO
# Provides:             logstash at startup, my_init from phusion. 
# Required-Start:       see /etc/init.d/logstash for more info 
# Required-Stop:        see /etc/init.d/logstash for more info 
# Short-Description:    You know, for search.
### END INIT INFO

set -e

#Source the lsb functions to perform the following operations
 . /lib/lsb/init-functions

log_daemon_msg "my_init is starting logstash"
/etc/init.d/logstash start

exit 0